# templates/duplicity-backup.service
[Unit]
Description=Run duplicity backup script

[Service]
Type=oneshot
ExecStart=/root/backup.sh
